# roles/node_exporter/files/etc/systemd/system/node_exporter.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/sbin/node_exporter --web.systemd-socket --collector.textfile.directory=/var/lib/node_exporter/ --collector.systemd
