# roles/automount/templates/share.mount.j2
[Unit]
Description={{ item.where }} mount

[Mount]
Where={{ item.where }}
What={{ item.what }}
Type={{ item.type }}
