# roles/automount/templates/share.automount.j2
[Unit]
Description={{ item.where }} automount

[Automount]
Where={{ item.where }}

[Install]
WantedBy=multi-user.target
