# roles/lego/files/lego@.timer
[Unit]
Description=Daily run of lego

[Timer]
OnBootSec=15m
OnUnitActiveSec=1d
Unit=lego@%i.service

[Install]
WantedBy=timers.target
