# roles/nullmailer/files/etc/cron.d/nullmailer-metrics
*/5 * * * * root /usr/local/bin/nullmailer_exporter > /var/lib/node_exporter/nullmailer_metrics.prom
