# roles/node_exporter/files/etc/cron.d/node-exporter-textfiles
0 * * * * node_exporter /usr/local/bin/repo_timestamp /var/lib/node_exporter/repo_timestamps.prom
