# roles/zfs/files/etc/cron.monthly/zfs-auto-snapshot
#!/bin/sh

exec zfs-auto-snapshot --default-exclude --quiet --label=monthly --keep=12 //
