# roles/zfs/files/etc/cron.weekly/zfs-auto-snapshot
#!/bin/sh

exec zfs-auto-snapshot --default-exclude --quiet --label=weekly --keep=8 //
