# roles/zfs/files/etc/cron.daily/zfs-auto-snapshot
#!/bin/sh

exec zfs-auto-snapshot --default-exclude --quiet --label=daily --keep=31 //
