# roles/zfs/files/etc/cron.d/zfs-auto-snapshot
PATH="/usr/bin:/bin:/usr/local/sbin:/usr/sbin:/sbin"

*/15 * * * * root zfs-auto-snapshot --default-exclude --quiet --label=frequent --keep=4 //
