# roles/cron/files/etc/cron.daily/journal
#!/usr/bin/env bash

journalctl -p 0..4 --since yesterday